# systemd/monitor-test.service
[Unit]
Description=Monitor the 'test' process and log restarts or monitoring issues
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitor-test.sh
User=root